# week_5/backup.service
[Unit]
Description=Script de Backup

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/auto-backup.sh