# week_5/backup.service
[Unit]
Description=Backup Script

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/auto-backup.sh