# week_2/nginx/test-nginx.timer
[Unit]
Description=Timer diari per comprovar Nginx

[Timer]
OnCalendar=daily
Unit=test-nginx.service

[Install]
WantedBy=timers.target