# nginx/test-nginx.service
[Unit]
Description=Comprovacio diaria de Nginx

[Service]
Type=oneshot
ExecStart=/usr/local/bin/test-nginx.sh